# almaany.service
[Unit]
Description=Almaany Crypto Trading Bot
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/root/almaany
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=append:/root/almaany/bot.log
StandardError=append:/root/almaany/bot.log

[Install]
WantedBy=multi-user.target[Unit]
Description=Almaany OKX demo trading bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/almaany
ExecStart=/path/to/almaany/venv/bin/python /path/to/almaany/main.py
Restart=always
RestartSec=10
EnvironmentFile=/path/to/almaany/.env
User=your_user_here

[Install]
WantedBy=multi-user.target
